# backup.service
[Unit]
Description=A script to backup my files every day 
[Service]
# The location of the  backup script
ExecStart=/home/elknawy/date.sh